# minecraft-bot.service.example
[Unit]
Description=Minecraft Gemini Bot
After=network.target

[Service]
Type=simple
User=your_user
Group=your_user_group
WorkingDirectory=/home/your_user/CraftBot
ExecStart=/usr/bin/node /home/your_user/CraftBot/craftbot.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
