# minecraft-bot.service.example
[Unit]
Description=CraftBot - Multi-Server Minecraft Chatbot with Gemini AI
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your_user
Group=your_user_group
WorkingDirectory=/home/your_user/CraftBot
ExecStart=/usr/bin/node /home/your_user/CraftBot/craftbot.js
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security and resource settings
PrivateTmp=true
NoNewPrivileges=true
LimitNOFILE=65536

# Environment
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=craftbot

[Install]
WantedBy=multi-user.targetion=Minecraft Gemini Bot
After=network.target